# packages/agents/steward/crontab
# Steward Crontab — DAO-Primary (VPS 3)
# All times in UTC

# Heartbeat check — every 5 minutes (CRITICAL)
*/5 * * * * /opt/cron/heartbeat-check.sh >> /var/log/agent/heartbeat-check.log 2>&1

# Proposal monitor — every 15 minutes (HIGH)
*/15 * * * * /opt/cron/proposal-monitor.sh >> /var/log/agent/proposal-monitor.log 2>&1

# Mod queue — every 30 minutes (MEDIUM)
*/30 * * * * /opt/cron/mod-queue.sh >> /var/log/agent/mod-queue.log 2>&1

# Dispute watcher — every hour at :35 (LOW)
35 * * * * /opt/cron/dispute-watcher.sh >> /var/log/agent/dispute-watcher.log 2>&1

# Stream claimer — every 4 hours at :10 (HIGH)
10 */4 * * * /opt/cron/stream-claimer.sh >> /var/log/agent/stream-claimer.log 2>&1

# Treasury health — every 6 hours at :55 (HIGH)
55 */6 * * * /opt/cron/treasury-health.sh >> /var/log/agent/treasury-health.log 2>&1

# Security audit — daily at 9:00 UTC (LOW)
0 9 * * * /opt/cron/security-audit.sh >> /var/log/agent/security-audit.log 2>&1

# Notification poll — every 5 minutes (HIGH)
*/5 * * * * /opt/cron/notification-poll.sh >> /var/log/agent/notification-poll.log 2>&1

# V3 cron jobs
15 */4 * * *  /opt/cron/reward-claimer.sh         # MEDIUM: every 4 hours at :15
0 */2 * * *   /opt/cron/loan-monitor.sh            # HIGH: every 2 hours
20 */4 * * *  /opt/cron/insurance-monitor.sh       # MEDIUM: every 4 hours at :20
*/15 * * * *  /opt/cron/lightning-watcher.sh        # HIGH: every 15 min

# Consensus resolver — every 2 minutes (HIGH)
*/2 * * * * /opt/cron/consensus-resolver.sh >> /var/log/agent/consensus-resolver.log 2>&1

# Anti-gaming monitoring (Steward = cross-agent monitor → all four at MEDIUM)
0 */12 * * *  /opt/cron/arb-health.sh >> /var/log/agent/arb-health.log 2>&1               # MEDIUM: every 12 hours
0 */12 * * *  /opt/cron/collusion-scan.sh >> /var/log/agent/collusion-scan.log 2>&1        # MEDIUM: every 12 hours
30 */12 * * * /opt/cron/mod-activity.sh >> /var/log/agent/mod-activity.log 2>&1            # MEDIUM: every 12 hours at :30
30 */12 * * * /opt/cron/vote-pattern.sh >> /var/log/agent/vote-pattern.log 2>&1            # MEDIUM: every 12 hours at :30
