# packages/agents/sentinel/crontab
# Sentinel Crontab — Moderator-Primary (VPS 1)
# All times in UTC

# Heartbeat check — every 5 minutes (CRITICAL)
*/5 * * * * /opt/cron/heartbeat-check.sh >> /var/log/agent/heartbeat-check.log 2>&1

# Mod queue — every 15 minutes (HIGH)
*/15 * * * * /opt/cron/mod-queue.sh >> /var/log/agent/mod-queue.log 2>&1

# Dispute watcher — every 30 minutes (MEDIUM)
*/30 * * * * /opt/cron/dispute-watcher.sh >> /var/log/agent/dispute-watcher.log 2>&1

# Proposal monitor — every hour at :20 (MEDIUM)
20 * * * * /opt/cron/proposal-monitor.sh >> /var/log/agent/proposal-monitor.log 2>&1

# Treasury health — every 4 hours at :45 (LOW)
45 */4 * * * /opt/cron/treasury-health.sh >> /var/log/agent/treasury-health.log 2>&1

# Security audit — daily at 9:00 UTC (LOW)
0 9 * * * /opt/cron/security-audit.sh >> /var/log/agent/security-audit.log 2>&1

# Notification poll — every 5 minutes (HIGH)
*/5 * * * * /opt/cron/notification-poll.sh >> /var/log/agent/notification-poll.log 2>&1

# V3 cron jobs
15 */4 * * *  /opt/cron/reward-claimer.sh         # MEDIUM: every 4 hours at :15
30 */6 * * *  /opt/cron/insurance-monitor.sh       # LOW: every 6 hours at :30

# Consensus resolver — every 2 minutes (HIGH)
*/2 * * * * /opt/cron/consensus-resolver.sh >> /var/log/agent/consensus-resolver.log 2>&1

# Anti-gaming monitoring (Sentinel = moderator-primary → mod-activity & collusion at HIGH)
0 */6 * * *   /opt/cron/collusion-scan.sh >> /var/log/agent/collusion-scan.log 2>&1       # HIGH: every 6 hours
0 */12 * * *  /opt/cron/mod-activity.sh >> /var/log/agent/mod-activity.log 2>&1            # HIGH: every 12 hours
0 3 * * 1     /opt/cron/arb-health.sh >> /var/log/agent/arb-health.log 2>&1               # MEDIUM: weekly Monday 03:00
0 4 * * 1     /opt/cron/vote-pattern.sh >> /var/log/agent/vote-pattern.log 2>&1            # MEDIUM: weekly Monday 04:00
