# packages/agents/sentinel/crontab
# Sentinel Crontab — Moderator-Primary (VPS 1)
# All times in UTC

# Heartbeat check — every 5 minutes (CRITICAL)
*/5 * * * * /opt/cron/heartbeat-check.sh >> /var/log/agent/heartbeat-check.log 2>&1

# Mod queue — every 15 minutes (HIGH)
*/15 * * * * /opt/cron/mod-queue.sh >> /var/log/agent/mod-queue.log 2>&1

# Dispute watcher — every 30 minutes (MEDIUM)
*/30 * * * * /opt/cron/dispute-watcher.sh >> /var/log/agent/dispute-watcher.log 2>&1

# Proposal monitor — every hour at :20 (MEDIUM)
20 * * * * /opt/cron/proposal-monitor.sh >> /var/log/agent/proposal-monitor.log 2>&1

# Treasury health — every 4 hours at :45 (LOW)
45 */4 * * * /opt/cron/treasury-health.sh >> /var/log/agent/treasury-health.log 2>&1
